# upload-pysqm-images.timer
[Unit]
Description=Periodically upload SQM images

[Timer]
OnCalendar=*:0/10
